# crontab
0 15 * * 6 cd /scripts/raspberry-scripts/ && python plan.py
0 8 * * * cd /scripts/raspberry-scripts/ && python notifyBirthday.py
0 4 * * * cd /scripts/raspberry-scripts/ && python syncBirthday.py